# certbot/crontab
0 7 * * 3 sh /certbot/run_certbot.sh 2>&1 /dev/stdout
