# certbot/crontab
0 7 * * 3 /bin/bash /certbot/run_certbot.sh 2>&1 /dev/stdout
